# deploy/agentd.service
[Unit]
Description=Agent OS daemon (agentd)
After=network-online.target
Wants=network-online.target
# Self-update probation (slice 5): if agentd crash-loops (5 restarts within 5 min),
# systemd stops thrashing + enters failed → OnFailure fires the probation rollback,
# which restores agentd.prev ONLY if a recent self-update confirmed (else it no-ops
# and leaves the failure to a human). Catches a LATENT crash the in-window watchdog
# can't (boots healthy → confirmed → dies minutes later).
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=apexos-rollback.service

[Service]
Type=simple
ExecStart=/usr/local/bin/agentd
Restart=on-failure
RestartSec=2

User=agentd
Group=agentd

# Root the daemon (and its file tools' relative paths) in APEX's workspace.
WorkingDirectory=/var/lib/agentd/workspace

Environment=AGENTD_PLUGINS_TOML=/etc/agentd/plugins.toml
Environment=AGENTD_POLICY_TOML=/etc/agentd/policy.toml
Environment=AGENTD_PARTS_INVENTORY=/etc/agentd/parts/inventory.toml
Environment=AGENTD_WORKSPACE=/var/lib/agentd/workspace
Environment=AGENTD_LOG=/var/lib/agentd/events
Environment=AGENTD_UI=/var/lib/agentd/ui
Environment=AGENTD_HARDWARE_WISHLIST=/var/lib/agentd/hardware-wishlist.md
# Where agentd drops USB eject-request files for the root drain service (it can't
# sudo under NoNewPrivileges) — under ReadWritePaths so the non-root daemon can write.
Environment=AGENTD_USB_EJECT_DIR=/var/lib/agentd/usb-eject
# Likewise for "Use this drive" prep requests (relabel/format a stick → exo-workspace).
Environment=AGENTD_USB_PREP_DIR=/var/lib/agentd/usb-prep
EnvironmentFile=-/etc/agentd/env

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/agentd /etc/agentd
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
