# deploy/apex-sensor-bridge.service
[Unit]
Description=ApexOS sensor bridge daemon
After=network.target agentd.service
Requires=agentd.service

[Service]
Type=simple
User=agentd
Group=agentd
EnvironmentFile=-/etc/agentd/env
Environment=SENSOR_INTERVAL_SECS=30
Environment=SENSOR_BRIDGE_HOST=localhost:8787
ExecStart=/usr/local/bin/apex-sensor-bridge
Restart=always
RestartSec=10

# To enable BME688 + MLX90640 readings, run a SensorHead dashboard (separate repo)
# and point the bridge at it — node-local, only sensor nodes have one:
#   Environment=SENSORHEAD_URL=http://localhost:8080   (or a systemd drop-in)
# Without it the bridge forwards CPU temp only (a missing dashboard never errors).

# Sandbox: runs as the agentd uid, so without this it could write everything agentd
# can. This binary opens NO device nodes — it reads /etc/hostname + /sys/class/thermal
# (CPU temp) and HTTP-polls the optional SensorHead dashboard above for BME688/MLX90640.
# The I2C hardware is read by that separate SensorHead service, NOT this binary — so
# PrivateDevices=true is correct (verified on a sensor node: readings still flow). It
# writes nothing to disk (no ReadWritePaths) and keeps /sys readable (no
# ProtectKernelTunables). The whole FS is otherwise read-only.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
