# deploy/apexos-rs-ui.service
[Unit]
Description=ApexOS-RS Native UI (Slint/KMS)
# Wants (NOT Requires) agentd: the UI rides through an agentd restart via its WS
# reconnect loop, so it must NOT be lifecycle-coupled. Requires= cascade-STOPS the UI
# whenever agentd stops — e.g. the self-update watchdog's `systemctl stop agentd`
# during a binary swap — and does NOT bring it back, leaving the kiosk dark until a
# manual restart. Wants keeps the boot ordering (After) without the teardown coupling.
After=agentd.service network.target
Wants=agentd.service

[Service]
User=root
Group=root
TTYPath=/dev/tty7
PAMName=login
StandardInput=tty
# Direct KMS/DRM rendering — no Wayland, no cage
Environment=SLINT_BACKEND=linuxkms
# agentd WS endpoint — override if port differs
Environment=AGENTD_WS=ws://localhost:8787/ws
EnvironmentFile=-/etc/agentd/env
ExecStart=/usr/local/bin/apexos-rs-ui
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=3
# Hardening — limit blast radius while still allowing KMS/DRM rendering
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
DevicePolicy=closed
DeviceAllow=/dev/dri/* rw
DeviceAllow=/dev/tty7 rw
DeviceAllow=/dev/input/* rw

[Install]
WantedBy=multi-user.target
